# systemd/system/servicemonitor.service
[Unit]
Description=Service Monitor Dashboard
After=network.target

[Service]
User=root
WorkingDirectory=/home/servicemonitor
ExecStart=/usr/bin/python3 /home/servicemonitor/app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
