# systemd/system/servicemonitor.service
[Unit]
Description=Service Monitor Dashboard
After=network.target

[Service]
User=root
WorkingDirectory=/home/servicemonitor
ExecStart=/home/servicemonitor/start_dashboard.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
